# deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:matching-engine]
command=/app/matching-engine/main
directory=/app/matching-engine
environment=CONFIG_PATH="/app/matching-engine/config.json",RUST_LOG="info"
autostart=true
autorestart=true
stdout_logfile=/var/log/matching-engine.log
stderr_logfile=/var/log/matching-engine-err.log
priority=10

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx-access.log
stderr_logfile=/var/log/nginx-error.log
priority=20

[program:admin-server]
command=python3 /app/admin/server.py
directory=/app/admin
autostart=true
autorestart=true
stdout_logfile=/var/log/admin-server.log
stderr_logfile=/var/log/admin-server-err.log
priority=30

[program:price-enforcer]
command=python3 /app/python_bots/MarketMaker.py
directory=/app/python_bots
autostart=false
autorestart=false
stdout_logfile=/var/log/price-enforcer.log
stderr_logfile=/var/log/price-enforcer-err.log
priority=40

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock
